# systemd/lxc-scaler.service
[Unit]
Description=LXC Auto-Scaler for Proxmox VE
After=network.target pve-cluster.service
Wants=pve-cluster.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/lxc-scaler/lxc_scaler.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
